# root/etc/services.d/02-notificator/run
#!/usr/bin/with-contenv

bash cd /app || exit

echo "[services.d] 02-notificator: Starting the app..."

exec \ s6-setuidgid abc python3 -u notificator.py