# root/etc/services.d/02-notificator/run
#!/usr/bin/with-contenv /bin/sh

echo "[services.d] 02-notificator: Starting the app..."

cd /app 

python notificator.py
