# root/etc/services.d/02-notificator/run
#!/usr/bin/with-contenv

cd /app

echo "[services.d] 02-notificator: Starting the app..."

s6-setuidgid abc python3 -u notificator.py