# root/etc/services.d/01-webhook/run
#!/usr/bin/with-contenv

cd /app

echo "[services.d] 01-webhook: Starting the webhooks..."

s6-setuidgid abc python3 -u print_urls.py
s6-setuidgid abc gunicorn -b :5445 -w 4 --access-logfile - main.app