# root/etc/services.d/01-webhook/run
#!/usr/bin/with-contenv

bash cd /app || exit

echo "[services.d] 01-webhook: Starting the webhooks..."

exec \ s6-setuidgid abc gunicorn -b :5445 -w 4 --access-logfile - main.app