# root/etc/services.d/01-webhook/run
#!/usr/bin/with-contenv /bin/sh

echo "[services.d] 01-webhook: Starting the webhooks..."

cd /app 
gunicorn -b :5445 -w 4 --access-logfile - main.app
