# root/etc/services.d/01-webhook/run
#!/usr/bin/with-contenv bash

cd /app || exit

exec \
	s6-setuidgid abc python3 -u print_urls.py

exec \
    s6-setuidgid abc waitress-serve --port=5445 main:app