# root/etc/services.d/01-webhook/run
#!/usr/bin/with-contenv bash

cd /app || exit

exec \
	s6-setuidgid abc python3 -u print_urls.py

exec \
    s6-setuidgid abc gunicorn -b :5445 -w 4 --access-logfile - main:app